# deploy/trux.service
[Unit]
Description=Trux backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dreamLand/jodulabs/trux
EnvironmentFile=%h/.trux/.env
ExecStart=/usr/bin/env pnpm --filter backend start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
