# deploy/trux.service.template
[Unit]
Description=Trux backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__TRUX_DIR__
EnvironmentFile=%h/.trux/.env
ExecStart=__PNPM__ --filter backend start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
